# systemd/mysql-proxy.service
[Unit]
Requires=mysql-proxy.service
After=mysql-proxy.service

[Service]
ExecStart=/lib/systemd/systemd-socket-proxyd 10.0.300.15:3306
